# install/aurora-web.service
[Unit]
Description=Aurora Web - LED Matrix Finger Paint Server
After=network.target

[Service]
Type=simple
User=debian
WorkingDirectory=/home/debian/AuroraV7
ExecStart=/home/debian/AuroraV7/.venv/bin/python -m uvicorn aurora_web.main:app --host 0.0.0.0 --port 80
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Give access to serial port
SupplementaryGroups=dialout

# Allow binding to privileged port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
